# deploy/sonos-lyrics.service
[Unit]
Description=Sonos Lyrics web app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/sonos-lyrics
ExecStart=/usr/bin/node server.js
Environment=NODE_ENV=production
Environment=PORT=3000
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
